# xrootd-proxy-renew/osg/xrootd-proxy-renew.cron
28 * * * * root /usr/bin/xrootd-gsi-renew --cron
